# oscam_update/timer_service/softcam-update.service
[Unit]
Description=update softcam.key

[Service]
Type=simple
ExecStart=/opt/oscam_update/softcam_update.sh

[Install]
WantedBy=multi-user.target